# systemd/systemd/user/wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnUnitActiveSec=15m
Unit=wallpaper.service

[Install]
WantedBy=default.target
